# python_sdk/bin/Linux64_x64/usbfs-memsize-systemd.service
[Unit]
Description=Set the size of the usbfs memory to 1000MB

[Service]
ExecStart=/bin/sh -c "echo 1000 > /sys/module/usbcore/parameters/usbfs_memory_mb"

[Install]
WantedBy=multi-user.target
